# .config/autostart/xbindkeys.desktop
[Desktop Entry]
Type=Application
Exec=xbindkeys
Hidden=false
X-GNOME-Autostart-enabled=true
Name=xbindkeys
Comment=Bind keys
